# deploy/systemd/openhalo-runtime.service.example
[Unit]
Description=OpenHalo Personal Runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openhalo
Group=openhalo
WorkingDirectory=/opt/openhalo
EnvironmentFile=/etc/openhalo/runtime.env
ExecStart=/opt/openhalo/.venv/bin/python -m personal_runtime.main \
  --host 127.0.0.1 \
  --port 8765 \
  --token-env OPENHALO_EDGE_TOKEN \
  --state-path ${OPENHALO_STATE_PATH} \
  --diagnostic-log-path ${OPENHALO_DIAGNOSTIC_LOG_PATH} \
  --runtime-config-path ${OPENHALO_RUNTIME_CONFIG_PATH}
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
